# roles/atuin/templates/atuin-sync.service.j2
[Unit]
Description=Atuin history sync

[Service]
Type=oneshot
ExecStart={{ atuin_bin_path }} sync
