# roles/atuin/files/atuin-sync.service
[Unit]
Description=Atuin history sync

[Service]
Type=oneshot
ExecStart=/usr/bin/atuin sync
